# board/batocera/allwinner/h700/rg40xx/fsoverlay/etc/udev/rules.d/80-disable-wlan0.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlan0", RUN+="/sbin/ip link set wlan0 down"
